# install-gitlab&drone/tools/crontab
0 6 * * 2-7 ubuntu /home/ubuntu/tools/backup.sh >/dev/null 2>&1